# roles/t460s/templates/udev_powermanagement.rules.j2
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/powermanagement false"
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/local/bin/powermanagement true"
